# crontab
0 4 * * * cd /usr/src/app && node dist/app.js || node dist/app.js
