# crontab
0 4 * * * cd /usr/src/app && npx ts-node --transpile-only src/app.ts || npx ts-node --transpile-only src/app.ts
