# crontab
0 4 * * * cd /usr/src/app && node --enable-source-maps dist/app.js || node --enable-source-maps dist/app.js
